# rules/ai_agent/ai_agent_openclaw_suspicious_file_write.yml
title: OpenClaw Suspicious File Write to Sensitive Paths
id: db43664e-2408-5e70-9b63-dcb0fabc8c5f
status: stable
description: |
  Detects OpenClaw file write or edit operations targeting sensitive system
  paths such as /etc/, shell startup files (.bashrc, .zshrc, .profile),
  crontab entries, and systemd service files. An AI agent writing to these
  locations may indicate persistence establishment or system tampering.
references:
  - https://attack.mitre.org/techniques/T1546/
  - https://attack.mitre.org/techniques/T1053/
author: AgentShield
date: "2026-02-07"
modified: "2026-02-24"
tags:
  - attack.persistence
  - attack.t1546
  - attack.t1053
logsource:
  product: ai_agent
  category: agent_events
  source: openclaw
detection:
  # Writes to /etc/
  selection_etc_write:
    source: openclaw
    event_type: file_write
    command|contains: '/etc/'
  selection_etc_edit:
    source: openclaw
    event_type: file_edit
    command|contains: '/etc/'
  # Shell startup file modifications
  selection_bashrc:
    source: openclaw
    event_type|contains:
      - file_write
      - file_edit
    command|contains:
      - '.bashrc'
      - '.bash_profile'
      - '.zshrc'
      - '.zprofile'
      - '.profile'
  # Crontab modifications via exec
  selection_crontab:
    source: openclaw
    event_type: tool_call
    command|contains: 'crontab'
  # Writes to cron directories
  selection_cron_write:
    source: openclaw
    event_type: file_write
    command|contains:
      - '/cron.d/'
      - '/crontab'
      - '/cron.daily/'
      - '/cron.hourly/'
  # systemd service file writes
  selection_systemd:
    source: openclaw
    event_type: file_write
    command|contains:
      - '/systemd/system/'
      - '/systemd/user/'
  # launchd plist writes
  selection_launchd:
    source: openclaw
    event_type: file_write
    command|contains:
      - '/LaunchAgents/'
      - '/LaunchDaemons/'
  # Authorized keys modification
  selection_authorized_keys:
    source: openclaw
    event_type|contains:
      - file_write
      - file_edit
    command|contains: '.ssh/authorized_keys'
  condition: selection_etc_write or selection_etc_edit or selection_bashrc or selection_crontab or selection_cron_write or selection_systemd or selection_launchd or selection_authorized_keys
falsepositives:
  - Legitimate system configuration management by authorised tools
  - Development scripts that modify shell profiles during environment setup
level: high
